# greenhouse.service
[Unit]
Description=Greenhouse Controller
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/greenhouse-controller
ExecStart=/usr/bin/python3 /home/pi/greenhouse-controller/backend/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
